# scripts/candy-machine.service
[Unit]
Description=Valentine's Candy Machine
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/vending
ExecStart=/home/pi/vending/venv/bin/python /home/pi/vending/backend/main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
